# deploy/systemd/sparsevllm-worker@.service
[Unit]
Description=Sparse-vLLM worker %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
UMask=0077
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.config/sparsevllm/worker-%i.env
ExecStart=/bin/bash -c 'cd "$SPARSEVLLM_REPO" && exec env PYTHONPATH="$SPARSEVLLM_REPO/src" "$SPARSEVLLM_PYTHON" -m sparsevllm.entrypoints.openai.api_server --model "$SPARSEVLLM_MODEL" --served-model-name "$SPARSEVLLM_SERVED_MODEL_NAME" --host 127.0.0.1 --port "$SPARSEVLLM_PORT" --engine-kwargs "$SPARSEVLLM_ENGINE_KWARGS" --request-log-dir "$SPARSEVLLM_REQUEST_LOG_DIR"'
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=default.target
